# raspool.service
[Unit]
Description=Pool Automation Controller

[Service]
ExecStart=/home/hburke/raspool/raspool.py
Restart=always
